# systemd/plowpilot-tensorrt.service
[Unit]
Description=PlowPilot AI-Vision TensorRT Engine Builder
Documentation=https://github.com/sompublic/AI-Vision
After=multi-user.target

[Service]
Type=oneshot
User=plowpilot
Group=plowpilot
WorkingDirectory=/app
# ✅ Safe configuration: FP16 only, no benchmark, with timeout
ExecStart=/app/scripts/build_engines.sh
# ✅ Add timeout protection (2 hours max)
TimeoutSec=7200
# ✅ Restart policy for failed builds
Restart=on-failure
RestartSec=300
StandardOutput=journal
StandardError=journal
SyslogIdentifier=plowpilot-tensorrt

[Install]
WantedBy=multi-user.target
